# etc/profile.d/metadose.sh
export RENADA_HOME=/srv/git/Renada-1
export CMP_LOGS=/var/log/cmp/system/logs
export TOMCAT_LOGS=/var/log/tomcat
export TOMCAT_WEBAPPS=/usr/share/tomcat/webapps
export PRAC_DEV=/srv/git/PracticeRepo-dev/mod-webapp-codebase